# server_setup/ansible/roles/rest_api/templates/atmrestapi.service.j2
[Unit]
Description=Service for API of ATM project
After=network.target

[Service]
Type=simple
User={{ service_user }}
WorkingDirectory={{ api_path }}
ExecStart=npm start
